# assembly/src/release/bin/solaris/activemq.xml
<!DOCTYPE service_bundle SYSTEM
"/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<!--
    Licensed to the Apache Software Foundation (ASF) under one or more
    contributor license agreements.  See the NOTICE file distributed with
    this work for additional information regarding copyright ownership.
    The ASF licenses this file to You under the Apache License, Version 2.0
    (the "License"); you may not use this file except in compliance with
    the License.  You may obtain a copy of the License at

    http://www.apache.org/licenses/LICENSE-2.0

    Unless required by applicable law or agreed to in writing, software
    distributed under the License is distributed on an "AS IS" BASIS,
    WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
    See the License for the specific language governing permissions and
    limitations under the License.
-->     
<service_bundle type='manifest' name='activemq'>

<service
        name='network/activemq'
        type='service'
        version='1'>

                <create_default_instance enabled='true' />

        <single_instance/>

        <dependency
                name='activemq-multi-user-server-milestone'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/milestone/multi-user-server' />
        </dependency>
	<method_context>
		<method_credential user='activemq' group='other'/>
	       <method_environment>
                     <envvar name="PATH" value="/usr/wsfm/bin:/usr/sfw/bin:/usr/ccs/bin:/usr/local/bin:/usr/bin:" />
                     <envvar name="PWD" value="/home/activemq/activemq" />
                     <envvar name="HOME" value="/home/activemq" />
                     <envvar name="ACTIVEMQ_BASE" value="/home/activemq/activemq" />
                     <envvar name="ACTIVEMQ_HOME" value="/home/activemq/activemq" />
               </method_environment>
	</method_context>


        <exec_method
                type='method'
                name='start'
                exec='/home/activemq/activemq/bin/activemq start'
                timeout_seconds='60' />

        <exec_method
        type='method'
        name='stop'
        exec=':kill'
        timeout_seconds='60' />
        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>Apache ActiveMQ
                        </loctext>
                </common_name>
		<documentation>
			<doc_link name='apache.org' uri='http://activemq.apache.org/' />
		</documentation>
        </template>
</service>

</service_bundle>